# debian/webengine.cron.d
# Clean up stale sessions
25      3       *       *       *       sjuser  /usr/share/webengine/app/webengine/manage.py cleanup
